# mycel-iso/airootfs/etc/sv/autologin/run
#!/bin/sh
exec chpst -u live env HOME=/home/live USER=live sway
